# systemd-user/openclaw-config-sync-push.path
[Unit]
Description=Watch Moltbot config for changes (sync + push)

[Path]
PathChanged=%h/.openclaw/openclaw.json

[Install]
WantedBy=default.target
