# systemd-user/openclaw-config-sync-push.service
[Unit]
Description=Sync Moltbot config to git repo and push

[Service]
Type=oneshot
WorkingDirectory=%h/openclaw-config
ExecStart=/usr/bin/env bash -lc 'bash %h/openclaw-config/sync-config-push.sh'
