# systemd-user/openclaw-config-sync.path
[Unit]
Description=Watch Moltbot config for changes

[Path]
PathChanged=%h/.openclaw/openclaw.json

[Install]
WantedBy=default.target
